# recipes-workaround/ocdm-workaround/files/wpeframework-ocdm-PREMIUMAPP-3272.service
[Unit]
Description=WPEFramework OCDM Initialiser workaround
Requires=wpeframework.service
After=wpeframework.service wpeframework-ocdm.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator OCDM
